# www/rpc/crontab
# Simple way to start the test thin server
@reboot /wikk/www/rpc/thin.sh
